# hello-api.service
[Unit]
Description=Hello API

[Service]
Type=simple
ExecStart=/home/ubuntu/bin/main
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
